# systemd_units/kafka.service
[Unit]
Description=Apache Kafka
Requires=zookeeper.service
After=zookeeper.service

[Service]
Type=simple

ExecStart=/local/s25staff/kafka/bin/kafka-server-start.sh %h/src/queueserver/kafka_server.properties
ExecStop=/local/s25staff/kafka/bin/kafka-server-stop.sh

[Install]
WantedBy=multi-user.target
